# src/roccat.service
[Unit]
Description=starts roccat-vulcan from bin dir

[Service]
ExecStart=/usr/bin/roccat-vulcan

[Install]
WantedBy=multi-user.target
